# cocktail-machine.service
[Unit]
Description=Cocktail Machine Tipsy Elite
After=network.target graphical-session.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/cocktail.bzh
Environment=DISPLAY=:0
Environment=PYTHONPATH=/home/pi/cocktail.bzh
ExecStartPre=/bin/sleep 10
ExecStart=/home/pi/cocktail.bzh/venv/bin/python /home/pi/cocktail.bzh/src/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# GPIO et hardware access
SupplementaryGroups=gpio spi i2c

# Variables d'environnement pour l'interface tactile
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_FBDEV=/dev/fb1

[Install]
WantedBy=multi-user.target